# ledmatrix.service
[Unit]
Description=LED Matrix Display Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/ledpi/LEDMatrix
ExecStart=/usr/bin/python3 /home/ledpi/LEDMatrix/display_controller.py
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=ledmatrix

[Install]
WantedBy=multi-user.target 